# v5.2/systemd/sypnose-coordinator.service
[Unit]
Description=SYPNOSE Coordinator v5.2
After=network.target

[Service]
Type=simple
User=gestoria
WorkingDirectory=/opt/sypnose
ExecStart=/usr/bin/node /opt/sypnose/bin/start.js
Restart=always
RestartSec=10
EnvironmentFile=/opt/sypnose/.env
Environment=NODE_ENV=production
LimitNOFILE=65536
MemoryMax=2G

[Install]
WantedBy=multi-user.target
